# upstart/lte.conf
# lte - LTE daemon
# automatically connects to lte

description	"LTE daemon"

start on runlevel [2345]
stop on runlevel [!2345]

respawn

exec /usr/sbin/lte-connect acs

